# service/openbsd
#!/bin/sh
#

daemon="/usr/bin/emt-chatserver"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
